# crontab
* * * * * root /root/twitternator.sh cron
