# crontab
* * * * * root $HOME/bin/twitternator.sh cron  >> /var/log/twitternator.log 2>&1
