# crontab
* * * * * twitternator /home/twitternator/twitternator.sh cron
