# crontab
* * * * * root /root/twitternator.sh cron  >> /var/log/twitternator.log 2>&1
